# blue-volume-discover/root/etc/s6/discover/run
#!/bin/bash

exec discosrv -limit-avg=10000 -limit-burst=1000000 -limit-cache=2046