# nginx/root/etc/s6/php-fpm/run
#!/bin/bash

if [ "$ALLOW_DEBUG" == "true" ]
then
 echo "Allowing for debug"
fi

exec php5-fpm -F